# sample.cron
CRON_TZ=America/New_York

* * * * * root log echo "1min beacon"
